# plexrpm.crontab
0 */2 * * *  root cd REPLACEWITHPATH  && ./checkplexpackage.pl 2>&1 >/root/plexupdate.log
